# docker/supervisord.conf
; =============================================================================
; SIAP - Sistem Informasi Akademik POLSA
; Supervisor Configuration — Full Stack Manager
; =============================================================================
;
; Process yang dikelola:
;   [1] frankenphp       — Web server (Caddy + PHP Worker Mode) → port 8000
;   [2] queue-default    — Queue Worker untuk job umum (Feeder Sync, KRS, dll)
;   [3] queue-whatsapp   — Queue Worker khusus antrian WA (anti-spam, rate-limited)
;   [4] laravel-scheduler — Cron scheduler internal Laravel
;   [5] laravel-reverb   — WebSocket server Reverb → port 9000 (internal)
;
; Hierarki Priority: frankenphp (10) > scheduler (15) > reverb (18) >
;                    queue-default (20) > queue-whatsapp (25)
; =============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; =============================================================================
; [1] FrankenPHP Web Server
; =============================================================================
; FrankenPHP menjalankan PHP dalam "worker mode" — PHP tidak dibunuh antar
; request, sehingga bootstrap Laravel hanya terjadi sekali per worker.
; Ini menghasilkan throughput yang jauh lebih tinggi dari PHP-FPM konvensional.
;
; SERVER_NAME diset ke ":8000" via ENV di Dockerfile/docker-compose.
; FrankenPHP membaca env SERVER_NAME secara otomatis.
; =============================================================================
[program:frankenphp]
command=frankenphp run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
startsecs=5
startretries=5
stopsignal=SIGTERM
stopwaitsecs=30
stderr_logfile=/var/log/supervisor/frankenphp.err.log
stdout_logfile=/var/log/supervisor/frankenphp.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
priority=10

; =============================================================================
; [2] Laravel Queue Worker — Default Queue
; =============================================================================
; Menangani job-job general:
;   - FeederSyncJob (sinkronisasi data Neofeeder — heavy, timeout 3600s)
;   - SendKrsNotificationJob (notifikasi approval KRS)
;   - KirimJadwalDosen / KirimJadwalPegawai
;   - KirimNotifikasiJadwal
;
; numprocs=2: 2 worker paralel untuk menangani antrian lebih cepat.
; timeout=3600: sesuai dengan $timeout di FeederSyncJob.
; max-time=14400: worker di-restart setiap 4 jam untuk mencegah memory leak.
; memory=512: batas memory per worker (MB) sesuai php.ini.
; =============================================================================
[program:queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work database --queue=default --sleep=3 --tries=3 --timeout=3600 --max-time=14400 --memory=512 --no-interaction
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=SIGTERM
stopwaitsecs=120
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=20

; =============================================================================
; [3] Laravel Queue Worker — WhatsApp Queue (Anti-Spam)
; =============================================================================
; Queue terpisah khusus untuk pengiriman pesan WhatsApp.
; Worker tunggal (numprocs=1) dengan rest=2 detik antar job untuk
; mencegah rate-limit/spam detection dari API WA.
;
; Menangani:
;   - SendWhatsappOtpJob (OTP login)
;   - SendWhatsappResetPasswordJob (reset password)
;
; timeout=60: pesan WA biasanya selesai dalam detik, 60s adalah batas aman.
; =============================================================================
[program:queue-whatsapp]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work database --queue=whatsapp --sleep=5 --tries=3 --timeout=60 --max-time=14400 --memory=128 --no-interaction
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=SIGTERM
stopwaitsecs=30
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-whatsapp.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
priority=25

; =============================================================================
; [4] Laravel Scheduler
; =============================================================================
; Menggantikan cron system-level — mengeksekusi command yang terdaftar di
; app/Console/Kernel.php (atau routes/console.php di Laravel 10).
; schedule:work melakukan polling setiap 60 detik (standar cron).
; =============================================================================
[program:laravel-scheduler]
command=php /app/artisan schedule:work --no-interaction
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=SIGTERM
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
priority=15

; =============================================================================
; [5] Laravel Reverb WebSocket Server
; =============================================================================
; Reverb menangani WebSocket untuk:
;   - FeederSyncProgress (real-time progress bar sync Neofeeder)
;   - Chat per jadwal kuliah (channel: chat.{jadwalId})
;   - Bimbingan akademik (channel: guidance.{studentId})
;   - Notifikasi (MessageSent, MessageRead events)
;
; Mendengarkan di port INTERNAL 9000 (non-default).
; Dari luar container, diakses via host port 8101 (lihat docker-compose.yml).
;
; CATATAN: Jika menggunakan docker-compose service `reverb` yang terpisah,
; nonaktifkan program ini dengan autostart=false untuk menghindari konflik.
; =============================================================================
[program:laravel-reverb]
command=php /app/artisan reverb:start --host=0.0.0.0 --port=9000 --no-interaction
autostart=false
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=SIGTERM
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/reverb.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
priority=18

; =============================================================================
; [GROUP] Semua worker queue dalam satu group untuk mudah di-manage
; =============================================================================
[group:queues]
programs=queue-default,queue-whatsapp
priority=999
