# docker/supervisord.conf
; =============================================================================
; SIAP - Sistem Informasi Akademik POLSA
; Supervisor Configuration
; =============================================================================
; Supervisor mengelola proses FrankenPHP dan Laravel Queue Worker.
; Jika salah satu proses mati, Supervisor akan otomatis me-restart-nya.
; =============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---------- FrankenPHP Web Server ----------
[program:frankenphp]
command=frankenphp run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/frankenphp-error.log
stdout_logfile=/var/log/supervisor/frankenphp-output.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=10

; ---------- Laravel Queue Worker (Default) ----------
[program:queue-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work database --queue=default --sleep=3 --tries=3 --max-time=7200 --timeout=600 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-worker-default.log
stdout_logfile_maxbytes=5MB
stopwaitsecs=30
priority=20

; ---------- Laravel Queue Worker (WhatsApp - Anti Spam) ----------
[program:queue-worker-whatsapp]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work database --queue=whatsapp --sleep=3 --rest=2 --tries=3 --max-time=7200 --timeout=600 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-worker-whatsapp.log
stdout_logfile_maxbytes=5MB
stopwaitsecs=30
priority=25

; ---------- Laravel Scheduler ----------
[program:laravel-scheduler]
command=php /app/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stdout_logfile_maxbytes=5MB
stopwaitsecs=30
priority=15

; ---------- Laravel Reverb (WebSockets) ----------
[program:laravel-reverb]
command=php /app/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-reverb.log
stdout_logfile_maxbytes=5MB
priority=15
